# ztp.conf
[program:ztp]
command=python3 /home/ztp/ztp/launch.py
numprocs=1
directory=/home/ztp/ztp/
user=root
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
